# root/etc/services.d/snapdrop-server/run
#!/usr/bin/with-contenv bash

cd /app/snapdrop/server || exit

exec \
	s6-setuidgid abc /usr/bin/node index.js
